# backend/fretquiz/fretquiz-backend.service
[Unit]
Description=FretQuiz backend service

[Service]
WorkingDirectory=/home/kdp/projects/pub/fretquiz/backend/fretquiz
